# tools/supervisor.conf
[program:grpc]
command = python3 -m konlpy_grpc server

[program:http]
command = ./konlpy_grpc_gateway -alsologtostderr -endpoint localhost:50051 -swagger_dir _generated
depends_on = grpc
